# sentry/files/supervisor.conf
{%- set service = salt['pillar.get']('supervisor:server:service:'+service_name) %}
{%- from "sentry/map.jinja" import server with context %}

{%- if service.name == 'web' %}
{%- set command = "bin/sentry start" %}
{%- elif service.name == 'worker' %}
{%- set command = "bin/sentry run worker" %}
{%- elif service.name == 'cron' %}
{%- set command = "bin/sentry run cron" %}
{%- endif %}

[program:{{ service_name }}]
directory={{ server.dir.base }}/run
environment=SENTRY_CONF="{{ server.dir.config }}"
command={{ server.dir.base }}/env/{{ command }}
user=sentry
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=syslog
stderr_logfile=syslog
